# infra/quadlet/caddy.container
# Quadlet unit for the Caddy reverse proxy.
# Install to: ~/.config/containers/systemd/caddy.container
#
# Caddy reads ~/Caddyfile (rewritten by the deploy script on each slot swap)
# and stores certificates in two named volumes managed by podman.

[Unit]
Description=Caddy reverse proxy
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/library/caddy:2-alpine
ContainerName=caddy
Network=host
Volume=%h/Caddyfile:/etc/caddy/Caddyfile:Z
Volume=caddy_data.volume:/data
Volume=caddy_config.volume:/config
Pull=missing

[Service]
Restart=on-failure
TimeoutStartSec=60

[Install]
WantedBy=default.target
