# infra/quadlet/kalandra-temporal-ui.container
# Temporal Web UI on :8233. The OCI security list keeps the port unreachable
# from outside — reach it through an SSH tunnel:
#   ssh -L 8233:localhost:8233 <vm> → http://localhost:8233

[Unit]
Description=Kalandra Temporal UI
After=network-online.target kalandra-temporal.service
Wants=network-online.target

[Container]
Image=docker.io/temporalio/ui:2.51.1
ContainerName=kalandra-temporal-ui
Network=host
Environment=TEMPORAL_ADDRESS=127.0.0.1:7233
Environment=TEMPORAL_UI_PORT=8233

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
