# infra/quadlet/kalandra-api-green.container
# Quadlet unit for the green slot of the API.
#
# At most one slot (blue OR green) should be enabled+running at a time.
# The CI/CD deploy script swaps slots on each deploy.
#
# Pull=missing means systemd uses whatever image is local at start time.
# The deploy script runs `podman pull` before restarting the service so
# the new image is always present when the unit starts.

[Unit]
Description=Kalandra API (green slot)
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/kalicz/demopage/api:latest
ContainerName=kalandra-api-green
Network=host
EnvironmentFile=%h/kalandra-api.env
Environment=ASPNETCORE_URLS=http://+:8081
Environment=ASPNETCORE_ENVIRONMENT=Production
Pull=missing

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
