# infra/quadlet/kalandra-mcp-green.container
# Quadlet unit for the green slot of the demopage MCP server (rootless, under the opc user).
# See kalandra-mcp-blue.container for the blue/green contract.

[Unit]
Description=Kalandra MCP server (green slot)
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/kalicz/demopage/mcp:latest
ContainerName=kalandra-mcp-green
Network=host
EnvironmentFile=%h/demoPage/kalandra-mcp.env
Environment=ASPNETCORE_HTTP_PORTS=8083
Environment=ASPNETCORE_ENVIRONMENT=Production
Pull=missing

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
