# infra/quadlet/kalandra-temporal-ui.container
# Temporal Web UI on :8233 — published as https://temporal.kalandra.tech through
# the shared Caddy with basic_auth (the UI has no authentication of its own).
# The OCI security list keeps :8233 itself unreachable from outside, so the
# authenticated Caddy vhost is the only way in.

[Unit]
Description=Kalandra Temporal UI
After=network-online.target kalandra-temporal.service
Wants=network-online.target

[Container]
Image=docker.io/temporalio/ui:2.51.1
ContainerName=kalandra-temporal-ui
Network=host
Environment=TEMPORAL_ADDRESS=127.0.0.1:7233
Environment=TEMPORAL_UI_PORT=8233

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
