# infra/host/reboot-notify.timer
# Checks every 15 minutes for a pending reboot. The script (reboot-notify.sh)
# throttles the actual Slack notifications to once per 24h, so a frequent check
# means fast detection without a noisy channel. Persistent=true catches up a
# missed run after the VM was off.
[Unit]
Description=Check every 15 min for a pending reboot

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
